# pi/etc/okavango.cron
0   */6 * * * root /sbin/reboot
*/5 *   * * * root /home/pi/okavango/pi/okavango
